# supervisor/haproxy.conf
[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg -p /var/run/haproxy.pid
numprocs=1
autostart=true
autorestart=false